# extra_scripts/coexistence.service
[Unit]
Description=firewall reset for docker coexistence

[Service]
Type=oneshot
ExecStart=/usr/local/bin/coexistence.sh